# pmtscope.service
[Unit]
Description=PMTscope - PMT Performance Parameter Diagnostic Tool
After=network.target

[Service]
Type=simple
User=yjj
WorkingDirectory=/home/yjj/pmtdatabase/pmt-scope
Environment="PATH=/home/yjj/pmtdatabase/pmt-scope/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PMTSCOPE_HOST=0.0.0.0"
Environment="PMTSCOPE_PORT=8501"
ExecStart=/home/yjj/pmtdatabase/pmt-scope/.venv/bin/streamlit run /home/yjj/pmtdatabase/pmt-scope/app.py \
    --server.address 0.0.0.0 \
    --server.port 8501 \
    --server.headless true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
